# lg_start
#!/sbin/openrc-run
# Looking Glass OpenRC service
touch /dev/shm/looking-glass && chown $(whoami):kvm /dev/shm/looking-glass && chmod 660 /dev/shm/looking-glass"
exit
