# lg_start
#!/sbin/openrc-run
# Looking Glass OpenRC service
/bin/sh /usr/bin/lg_start.sh
exit
